# dashboard/rootfs/etc/s6-overlay/s6-rc.d/svc-chromium/run
#!/command/with-contenv sh

export DISPLAY=:0
echo "[svc-chromium] Starting Chromium service on DISPLAY=$DISPLAY"

log_xset_state() {
    LABEL="$1"
    echo "[svc-chromium] X11 display power state ($LABEL):"
    xset q 2>&1 | while IFS= read -r line; do
        echo "[svc-chromium]   $line"
    done
}

run_xset_guard() {
    DESC="$1"
    shift
    echo "[svc-chromium] Applying display guard ($DESC): xset $*"
    OUTPUT="$(xset "$@" 2>&1)"
    STATUS="$?"

    if [ "$STATUS" -ne 0 ]; then
        echo "[svc-chromium] WARNING: xset $* failed with status $STATUS"
    fi

    if [ -n "$OUTPUT" ]; then
        echo "$OUTPUT" | while IFS= read -r line; do
            echo "[svc-chromium]   $line"
        done
    fi

    log_xset_state "after xset $*"
}

monitor_xset_state() {
    while true; do
        sleep 300
        log_xset_state "periodic 5min check"
    done
}

# Poll until Xorg is accepting connections (up to 60s)
MAX=120
CT=0
until xdpyinfo >/dev/null 2>&1; do
    sleep 0.5
    CT=$(( CT + 1 ))
    if [ "$CT" -ge "$MAX" ]; then
        echo "[svc-chromium] FATAL: Xorg never became ready on :0"
        exit 1
    fi
done
echo "[svc-chromium] Xorg ready"
log_xset_state "before kiosk blanking guard"

# Disable X screen saver blanking and DPMS at runtime.
# The X server's built-in screen-saver timer defaults to 600 s. That can blank
# the framebuffer after 10 min of no input, causing the monitor to report
# "no signal". Keep the DPMS extension available so the dashboard can later
# force the display off/on when presence-based power control is enabled.
run_xset_guard "set screen saver timeout to zero" s 0 0
run_xset_guard "disable screen saver timer" s off
run_xset_guard "disable screen saver blanking" s noblank
run_xset_guard "set DPMS timeouts to zero" dpms 0 0 0
run_xset_guard "disable automatic DPMS" -dpms
log_xset_state "after kiosk blanking guard"
monitor_xset_state &
MONITOR_PID="$!"

# Poll until Python server is accepting connections (up to 60s)
CT=0
until python3 -c 'import urllib.request; urllib.request.urlopen("http://127.0.0.1:8080/config", timeout=1).read()' >/dev/null 2>&1; do
    sleep 0.5
    CT=$(( CT + 1 ))
    if [ "$CT" -ge "$MAX" ]; then
        echo "[svc-chromium] FATAL: Python server never started"
        exit 1
    fi
done
echo "[svc-chromium] Python server ready"

# Hide mouse cursor after 1 second of inactivity (no pointer on this kiosk)
unclutter -idle 1 -root &

# Chromium requires a D-Bus session bus
eval $(dbus-launch --sh-syntax)
export DBUS_SESSION_BUS_ADDRESS

echo "[svc-chromium] Launching Chromium kiosk"
chromium \
    --no-sandbox \
    --disable-dev-shm-usage \
    --kiosk \
    --window-size=1920,1080 \
    --window-position=0,0 \
    --ozone-platform=x11 \
    --no-first-run \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --disable-restore-session-state \
    --noerrdialogs \
    --disable-background-networking \
    --disable-component-update \
    --disable-default-apps \
    --disable-domain-reliability \
    --disable-extensions \
    --disable-sync \
    --disable-translate \
    --disable-features=TranslateUI,MediaRouter,OptimizationHints,PushMessaging,Vulkan,DefaultANGLEVulkan,VulkanFromANGLE,WebGPU,DawnGraphite,OptimizationGuideModelDownloading,OnDeviceModelExecution \
    --overscroll-history-navigation=0 \
    --disable-gpu \
    --disable-gpu-compositing \
    --disable-vulkan \
    --disable-accelerated-2d-canvas \
    --disable-accelerated-video-decode \
    --disable-client-side-phishing-detection \
    --safebrowsing-disable-auto-update \
    --check-for-update-interval=31536000 \
    --user-data-dir=/tmp/chromium-kiosk \
    --app=http://localhost:8080 \
    2>&1

STATUS=$?
kill "$MONITOR_PID" 2>/dev/null || true
echo "[svc-chromium] Chromium exited with status $STATUS"
exit "$STATUS"
