# dashboard/rootfs/etc/s6-overlay/s6-rc.d/svc-chromium/run
#!/command/with-contenv sh

export DISPLAY=:0
echo "[svc-chromium] Starting Chromium service on DISPLAY=$DISPLAY"

# Poll until Xorg is accepting connections (up to 60s)
MAX=120
CT=0
until xdpyinfo >/dev/null 2>&1; do
    sleep 0.5
    CT=$(( CT + 1 ))
    if [ "$CT" -ge "$MAX" ]; then
        echo "[svc-chromium] FATAL: Xorg never became ready on :0"
        exit 1
    fi
done
echo "[svc-chromium] Xorg ready"

# Disable X screen saver blanking and DPMS at runtime.
# xorg.conf disables the DPMS protocol extension, but the X server's built-in
# screen-saver timer (default 600 s) is independent — it blanks the framebuffer
# after 10 min of no input, causing the monitor to report "no signal".
xset s off
xset -dpms
xset s noblank

# Poll until Python server is accepting connections (up to 60s)
CT=0
until python3 -c 'import urllib.request; urllib.request.urlopen("http://127.0.0.1:8080/config", timeout=1).read()' >/dev/null 2>&1; do
    sleep 0.5
    CT=$(( CT + 1 ))
    if [ "$CT" -ge "$MAX" ]; then
        echo "[svc-chromium] FATAL: Python server never started"
        exit 1
    fi
done
echo "[svc-chromium] Python server ready"

# Hide mouse cursor after 1 second of inactivity (no pointer on this kiosk)
unclutter -idle 1 -root &

# Chromium requires a D-Bus session bus
eval $(dbus-launch --sh-syntax)
export DBUS_SESSION_BUS_ADDRESS

echo "[svc-chromium] Launching Chromium kiosk"
chromium \
    --no-sandbox \
    --disable-dev-shm-usage \
    --kiosk \
    --ozone-platform=x11 \
    --no-first-run \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --disable-restore-session-state \
    --noerrdialogs \
    --disable-background-networking \
    --disable-component-update \
    --disable-default-apps \
    --disable-domain-reliability \
    --disable-extensions \
    --disable-sync \
    --disable-translate \
    --disable-features=TranslateUI,MediaRouter,OptimizationHints,PushMessaging,Vulkan,DefaultANGLEVulkan,VulkanFromANGLE,WebGPU,DawnGraphite,OptimizationGuideModelDownloading,OnDeviceModelExecution \
    --overscroll-history-navigation=0 \
    --disable-gpu \
    --disable-gpu-compositing \
    --disable-vulkan \
    --disable-accelerated-2d-canvas \
    --disable-accelerated-video-decode \
    --disable-client-side-phishing-detection \
    --safebrowsing-disable-auto-update \
    --check-for-update-interval=31536000 \
    --window-position=0,0 \
    --window-size=1920,1080 \
    --user-data-dir=/tmp/chromium-kiosk \
    --app=http://localhost:8080 \
    2>&1

STATUS=$?
echo "[svc-chromium] Chromium exited with status $STATUS"
exit "$STATUS"
