# dashboard/rootfs/etc/s6-overlay/s6-rc.d/svc-python/run
#!/command/with-contenv sh
exec python3 /app/server.py 2>&1
